# deploy/kurisu-web.service
# Template systemd unit for running kurisu-web in the background. See
# README.md's deployment section for setup instructions.
[Unit]
Description=Kurisu Web (stats frontend)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/path/to/kurisu-web
ExecStart=/path/to/kurisu-web/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8081 --proxy-headers --forwarded-allow-ips 127.0.0.1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
