# deploy/kurisu-web.service
# Template systemd unit for running kurisu-web in the background. See
# README.md's deployment section for setup instructions.
[Unit]
Description=Kurisu Web (stats frontend)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/path/to/kurisu-web
# Host/port + proxy-header trust come from WEB_BIND / WEB_FORWARDED_ALLOW_IPS in .env (see run.py).
ExecStart=/path/to/kurisu-web/.venv/bin/python run.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
